# deployment/local/openwork-team.service.example
[Unit]
Description=OpenWork Team Host (LAN)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=CHANGE_ME_USER
Group=CHANGE_ME_GROUP
WorkingDirectory=/opt/openwork/openwork
EnvironmentFile=/etc/openwork/openwork.env
ExecStart=/usr/bin/env bash -lc 'pnpm --filter openwork-orchestrator dev -- serve --workspace "$WORKSPACE_DIR" --remote-access --openwork-host 0.0.0.0 --openwork-port "$OPENWORK_PORT" --connect-host "$OPENWORK_CONNECT_HOST" --openwork-token "$OPENWORK_TOKEN" --openwork-host-token "$OPENWORK_HOST_TOKEN" --approval manual'
Restart=always
RestartSec=3
TimeoutStopSec=25

[Install]
WantedBy=multi-user.target

